# .crontab
0 * * * * python /home/vagrant/sandbox/scraper/scraper.py
